# package/wattlined/etc/hotplug.d/iface/95-wattline
#!/bin/sh

case "${ACTION:-}" in
	ifup|ifupdate)
		/usr/lib/wattline/firewall-sync
		;;
esac
