# package/wattlined/etc/init.d/wattlined
#!/bin/sh /etc/rc.common
# procd service for wattlined
USE_PROCD=1
START=95
STOP=10
RUNTIME_DIR=/var/run/wattline
LISTENER_STATE="$RUNTIME_DIR/listener-state"

listener_state() {
	for key in http_enabled http_addr4 http_addr6 port https_enabled https_addr4 https_addr6 https_port tls_cert tls_key; do
		printf '%s=%s\n' "$key" "$(uci -q get "wattline.main.$key" 2>/dev/null || true)"
	done
}

save_listener_state() {
	mkdir -p "$RUNTIME_DIR"
	listener_state > "$LISTENER_STATE.new"
	mv "$LISTENER_STATE.new" "$LISTENER_STATE"
}

start_service() {
	/usr/lib/wattline/firewall-sync
	save_listener_state
	procd_open_instance
	procd_set_param command /usr/bin/wattlined -config /etc/config/wattline
	procd_set_param respawn 3600 5 5
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}

reload_service() {
	/usr/lib/wattline/firewall-sync
	current="$(listener_state)"
	previous="$(cat "$LISTENER_STATE" 2>/dev/null || true)"
	if [ "$current" != "$previous" ]; then
		/etc/init.d/wattlined restart
		return
	fi
	# Rule-only changes retain the BLE link and reload with SIGHUP.
	local pid
	for pid in $(pgrep -f '/usr/bin/wattlined -config /etc/config/wattline' || true); do
		kill -HUP "$pid"
	done
}

service_triggers() {
	procd_add_reload_trigger "wattline"
}
