# package/wattlined/etc/init.d/wattlined
#!/bin/sh /etc/rc.common
# procd service for wattlined
USE_PROCD=1
START=95
STOP=10
RUNTIME_DIR="${RUNTIME_DIR:-/var/run/wattline}"
DAEMON_STATE="${DAEMON_STATE:-$RUNTIME_DIR/daemon-state}"
FIREWALL_SYNC="${FIREWALL_SYNC:-/usr/lib/wattline/firewall-sync}"
SERVICE_SCRIPT="${SERVICE_SCRIPT:-/etc/init.d/wattlined}"
PGREP="${PGREP:-pgrep}"

daemon_state() {
	# wattlined only reloads rule sections on SIGHUP. Any main-section change
	# therefore needs a process restart, including auth, discovery, and WAN
	# policy in addition to listener/TLS settings.
	uci -q show wattline.main 2>/dev/null || true
}

save_daemon_state() {
	mkdir -p "$RUNTIME_DIR"
	chmod 0700 "$RUNTIME_DIR"
	daemon_state > "$DAEMON_STATE.new"
	chmod 0600 "$DAEMON_STATE.new"
	mv "$DAEMON_STATE.new" "$DAEMON_STATE"
}

start_service() {
	"$FIREWALL_SYNC" || return 1
	save_daemon_state || return 1
	procd_open_instance
	procd_set_param command /usr/bin/wattlined -config /etc/config/wattline
	procd_set_param respawn 3600 5 5
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}

reload_service() {
	"$FIREWALL_SYNC" || return 1
	current="$(daemon_state)"
	previous="$(cat "$DAEMON_STATE" 2>/dev/null || true)"
	if [ "$current" != "$previous" ]; then
		"$SERVICE_SCRIPT" restart
		return
	fi
	# Rule-only changes retain the BLE link and reload with SIGHUP.
	local pid
	for pid in $("$PGREP" -f '/usr/bin/wattlined -config /etc/config/wattline' || true); do
		kill -HUP "$pid"
	done
}

service_triggers() {
	procd_add_reload_trigger "wattline"
}
