# package/wattline-rtl8761b/etc/init.d/wattline-rtl8761b
#!/bin/sh /etc/rc.common

START=15
STOP=85

DRIVERCTL="${DRIVERCTL:-/usr/lib/wattline/rtl8761b/driverctl}"

start() {
	"$DRIVERCTL" activate
}

stop() {
	return 0
}
